# fpms.service.in
[Unit]
Description=FPMS daemon responsible for OLED management

[Service]
ExecStart=@bindir@/oled-start
WorkingDirectory=@datadir@/fpms

[Install]
WantedBy=multi-user.target
